# deploy/systemd/ztool-license-server.service
[Unit]
Description=ZTool License TCP Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ztool
Group=ztool
WorkingDirectory=/opt/ztool-license-server/app/license-server
EnvironmentFile=/etc/ztool-license-server/ztool-license-server.env
ExecStart=/opt/ztool-license-server/.venv/bin/python -m ztool_license_server.server
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ztool-license-server /var/log/ztool-license-server

[Install]
WantedBy=multi-user.target
